# detector/gewu.service
[Unit]
Description=格物 relay verification web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/gewu-detector
Environment=PYTHONUNBUFFERED=1
Environment=GEWU_JOBS_DIR=/opt/gewu-detector/web_data/jobs
Environment=GEWU_WISHLIST_PATH=/opt/gewu-detector/web_data/wishlist.txt
EnvironmentFile=-/etc/gewu.env
ExecStartPre=/usr/bin/test -d /opt/gewu-detector/web_data/jobs
ExecStartPre=/usr/bin/test -w /opt/gewu-detector/web_data/jobs
ExecStart=/opt/gewu-detector/venv/bin/python -m uvicorn web.server:app \
  --host 127.0.0.1 --port 8765 \
  --workers 1 \
  --log-level info \
  --proxy-headers \
  --forwarded-allow-ips='127.0.0.1'
Restart=on-failure
RestartSec=3
TimeoutStopSec=45
StandardOutput=journal
StandardError=journal
User=gewu
Group=gewu
UMask=0077
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/gewu-detector/web_data
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
