# detector/gewu-monitor@.timer
[Unit]
Description=Run 格物 relay monitor %i every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=10min
Persistent=true
Unit=gewu-monitor@%i.service

[Install]
WantedBy=timers.target
