# detector/suyuan-monitor@.service
[Unit]
Description=溯源 scheduled relay monitor (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/suyuan-monitor/%i.env
ConditionPathExists=/etc/suyuan-monitor/%i.key

[Service]
Type=oneshot
WorkingDirectory=/opt/suyuan-detector
Environment=PYTHONUNBUFFERED=1
Environment=SUYUAN_JOBS_DIR=/opt/suyuan-detector/web_data/jobs
EnvironmentFile=/etc/suyuan-monitor/%i.env
LoadCredential=api_key:/etc/suyuan-monitor/%i.key
ExecStart=/opt/suyuan-detector/venv/bin/suyuan monitor-once --target-id %i
User=suyuan
Group=suyuan
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/suyuan-detector/web_data
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
