# detector/suyuan-monitor@.timer
[Unit]
Description=Run 溯源 relay monitor %i every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=10min
Persistent=true
Unit=suyuan-monitor@%i.service

[Install]
WantedBy=timers.target
